# etc/rc.local
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.

TMPDIR="/run/shm/honza/tmp"
mkdir -p $TMPDIR

chown -R honza.honza $TMPDIR
chmod -R 700 $TMPDIR
ln -s -f $TMPDIR /home/honza/.ramdisk

chgrp adm /dev/input/event*

exit 0
